# resources/systemd/arinna.target
[Unit]
Description=Arinna service group

[Install]
WantedBy=default.target
