# resources/systemd/arinna-load-balancer.timer
[Unit]
Description=Arinna load balancer timer
PartOf=arinna.target

[Timer]
OnUnitActiveSec=10s
AccuracySec=1s

[Install]
WantedBy=arinna.target
