# resources/systemd/arinna.target
[Unit]
Description=Arinna service group

[Install]
WantedBy=multi-user.target
